# systemd/kano-common-boards.service
#
# kano-common-boards
#

[Unit]
Description=Kano Boards Daemon
IgnoreOnIsolate=true

[Service]
ExecStart=/usr/bin/kano-boards-daemon
